# usingpc/usingpc.timer
[Unit]
Description=Run usingpc timer

[Timer]
OnCalendar=*:0/2
Persistent=true     
 
[Install]
WantedBy=timers.target
